# dogwatch.service
[Unit]
Description=DogWatch monitoring daemon
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/dogwatch/dogwatch.sh daemon
Restart=always

[Install]
WantedBy=multi-user.target
